# services/inotify/run
#!/bin/sh

/sbin/setuser abc inotifywait -m -r -e modify /config/apache --format %w%f | while read file; do
  if [ "${file##*.}" = "conf" ]; then
    sv down /etc/service/web/
    sv up /etc/service/web
    echo "$file changed, restarting"
  fi
done
