# services/inotify/run
#!/bin/sh

while /sbin/setuser abc inotifywait -e modify /config/apache/site-confs/; do
  sv down /etc/service/web/
  sv up /etc/service/web/
  echo "File in Config Folder Changed, Restarted"

done
