# services/reverseproxy/run
#!/bin/bash

exec /sbin/setuser abc /usr/sbin/apache2ctl -D FOREGROUND
